# launchd/com.dlowe.flicks.publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Daily build + publish of flicks. Installed by ./install-launchagent.sh, which
     substitutes __FLICKS_DIR__ for this repo's path. A LaunchAgent (per-user) so it
     runs in your logged-in session with your network + SSH credentials. Unlike cron,
     launchd runs a run missed during sleep when the machine next wakes. -->
<plist version="1.0"><dict>
  <key>Label</key><string>com.dlowe.flicks.publish</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>__FLICKS_DIR__/publish.sh</string>
  </array>
  <key>WorkingDirectory</key><string>__FLICKS_DIR__</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <!-- Stream Python's output to the log instead of block-buffering it (so a long
       run shows progress, not an empty file). -->
  <key>EnvironmentVariables</key>
  <dict><key>PYTHONUNBUFFERED</key><string>1</string></dict>
  <key>StandardOutPath</key><string>__FLICKS_DIR__/.publish.log</string>
  <key>StandardErrorPath</key><string>__FLICKS_DIR__/.publish.log</string>
  <!-- No ThrottleInterval: this is a once-a-day job, not a respawn loop, and the
       300s throttle just parks manual re-runs in "spawn scheduled" (looks hung). -->
</dict></plist>
